# Before-After/custom.target
[Unit]
Description=Custom Target
Documentation=https://github.com/nmasse-itix/Systemd-Examples
AllowIsolate=no
